# config/docker/supervisord-web.conf
[supervisord]
nodaemon=true
user=root

[program:web]
command=/home/django_admin/website_api/venv/bin/gunicorn website_api.wsgi:application --bind 0.0.0.0:8000
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/web.err.log
stdout_logfile=/var/log/web.out.log
